# dot_local/etc/sh/profile.d/60-editor.sh
# Text editors, pagers, etc.

export EDITOR="/opt/homebrew/bin/nvim"
export VISUAL="/opt/homebrew/bin/nvim"
export PAGER="/usr/bin/less"
