# dot_local/etc/sh/profile.d/90-homebrew.sh
# Homebrew
eval "$(/opt/homebrew/bin/brew shellenv)"

# Homebrew Github Token
source $XDG_LIB_HOME/private/homebrew/homebrew-github-token.sh

